# deploy/always-on/systemd/redeploy-website.service
[Unit]
Description=reDeploy marketing website (always-on)
After=network.target
# Ordering only (NOT a Requires/Wants dependency): start after redeploy-app so
# their ExecStartPre `pnpm install` runs don't race on the shared .serve/active
# workspace node_modules on a cold boot. The website still starts fine on its own
# if redeploy-app is stopped/failed — After= affects ordering, not requirement.
After=redeploy-app.service

[Service]
# Rides the SAME repo-local `.serve/active` symlink as redeploy-app.service —
# so /test-pr --serve/--reset repoint the served checkout for the website too.
WorkingDirectory=%h/Development/thesolidchain/reDeploy/.serve/active

# nvm-managed node/pnpm + foundry are NOT on the default systemd user-service PATH;
# pin them explicitly. Replace <NODE_BIN_DIR> with the output of:
#   dirname "$(which node)"
# If you switch node versions via nvm later, update this path too.
Environment=PATH=<NODE_BIN_DIR>:%h/.local/share/pnpm:%h/.foundry/bin:/usr/local/bin:/usr/bin:/bin

# No EnvironmentFile= here, unlike redeploy-app.service — the website is a
# static marketing site with no backend, no RPC, no keys.

ExecStartPre=/bin/bash -c 'pnpm install --prefer-offline'
# website Vite dev server on :5180 (port pinned in apps/website/vite.config.ts)
ExecStart=/bin/bash -c 'pnpm -F @redeploy/website dev'
Restart=always
RestartSec=5

# The ExecStartPre install above can take well over systemd's 90s default
# TimeoutStartSec on a cold checkout or crash recovery (empty pnpm store),
# which would otherwise get the unit SIGTERM'd mid-install and eventually
# left failed (StartLimitBurst). Disable the start timeout so slow first-time /
# cold-start installs always get to finish. See docs/ALWAYS-ON-TUNNEL.md section 10.
TimeoutStartSec=0

[Install]
WantedBy=default.target
