# .claude/systemd/claude-rc.service
# @orchestrator-managed claude-rc-service v4
# systemd (user) unit TEMPLATE (issue #102). Installed into
# ~/.config/systemd/user/claude-rc-<repo-slug>.service by
# `.claude/scripts/arm-loop.sh`, which substitutes the __PLACEHOLDER__ tokens
# below and runs `systemctl --user enable --now`.
#
# `claude remote-control` has no documented headless mode, so this unit runs
# it inside a DETACHED tmux session (`tmux new -d -s rc-<repo-slug>`): tmux
# gives systemd a stable child to track AND gives a human a local attach point
# (`tmux attach -t rc-<repo-slug>`) to see the QR/status or restart it by hand.
#
# Type=oneshot + RemainAfterExit=yes (not Type=forking): `tmux new -d` talks
# to the tmux SERVER over a socket and exits immediately once the detached
# session exists — if a tmux server is already running, ExecStart's own PID
# has no parent/child relationship to the long-lived process at all, so
# systemd cannot reliably track it as a "forked" child. This is the standard
# pattern for supervising a tmux/screen-managed daemon from systemd. Caveat:
# systemd only observes ExecStart's (successful) exit, not the health of the
# `claude remote-control` process running inside the tmux session — if THAT
# process itself crashes, tmux keeps the (now-empty) session and systemd sees
# nothing wrong. Inspect with `tmux attach -t rc-<repo-slug>`; restart with
# `systemctl --user restart claude-rc-<repo-slug>.service`.
#
# DO NOT hand-edit the INSTALLED copy under ~/.config/systemd/user/ — it will
# be silently overwritten the next time arm-loop.sh runs. Edit THIS checked-in
# template instead (then re-run arm-loop.sh), or fork it under a different
# name if you want a permanently custom copy. Re-run `/orchestrator:sync` to
# pick up plugin updates to this template before re-arming.
[Unit]
Description=claude remote-control server for __REPO_SLUG__ (detached tmux)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=__WORKDIR__
ExecStartPre=-/usr/bin/tmux kill-session -t rc-__REPO_SLUG__
ExecStart=/usr/bin/tmux new -d -s rc-__REPO_SLUG__ -c __WORKDIR__ /bin/bash -c 'export PATH=__CLAUDE_DIR__:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin; exec __CLAUDE_BIN__ remote-control --name __RC_NAME__ --remote-control-session-name-prefix __REPO_SLUG__ --capacity __CAPACITY__ --permission-mode __PERMISSION_MODE__'
ExecStop=-/usr/bin/tmux kill-session -t rc-__REPO_SLUG__

[Install]
WantedBy=default.target
